# preconf/smrd.service.example
[Unit]
Description=SMR Peer-to-Peer Network Daemon

[Service]
Type=simple
User=smr
ExecStart=/usr/bin/smrd --conf=/etc/opt/smr/smrd.cfg

[Install]
WantedBy=multi-user.target
